# etc/sentinel.timer
[Unit]
Description=Run sentinel.service every 2 minutes

[Timer]
OnCalendar=*:0/2
