# roles/sentinel/files/sentinel-upload.path
[Path]
PathChanged=/tmp/sentinel.jpg
Unit=sentinel-upload.service

[Install]
WantedBy=multi-user.target
