# roles/sentinel/files/sentinel.service
[Unit]
Description=Renders camera output to a file
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sentinel.sh
